# babel-rulebook.service
[Unit]
Description=Babel Rulebook Static Server
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/babel-rulebook
ExecStart=/home/exedev/.bun/bin/bun run server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
